# packaging/linux/cannonball.service
[Unit]
Description=Cannonball corporate messaging platform
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cannonball
Group=cannonball
WorkingDirectory=/opt/cannonball
EnvironmentFile=/etc/cannonball/cannonball.env
ExecStart=/opt/cannonball/bin/cannonball
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
